# vendor/github.com/vmware/vic/isos/appliance/permissions.service
[Unit]
Description=Modify file permissions for VIC target
After=basic.target

[Service]
Type=oneshot
ExecStart=/etc/systemd/scripts/permissions-setup
RemainAfterExit=yes

[Install]
WantedBy=vic.target