# vendor/github.com/vmware/vic/isos/appliance/vic-init.service
[Unit]
Description=Init process for VIC components
After=permissions.service

[Service]
Type=idle

LimitAS=infinity
LimitRSS=infinity
LimitCORE=infinity
LimitNOFILE=65536

# This is not configured for restart as currently vic-init isn't able to
# reacquire knowledge about whether sessions are running.
# Instead there is a top level trigger to restart the VM - reverting to
# known good VM state
ExecStart=/sbin/vic-init

[Install]
WantedBy=vic.target
Wants=basic.target
Wants=multi-user.target
