# vendor/github.com/vmware/vic/isos/bootstrap/tether.service
[Unit]
Description=Tether
After=basic.target

[Service]
Type=idle
ExecStart=/bin/sh -c '/bin/bootstrap >/dev/ttyS1 2>&1'

[Install]
WantedBy=vic.target
