# vendor/github.com/vmware/vic/isos/bootstrap/vic.target
[Unit]
Description=vSphere Integrated Containers
Requires=basic.target
After=basic.target
AllowIsolate=yes